# scripts/computenook.service
[Unit]
Description=算力小筑 Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/computenook
EnvironmentFile=-/opt/computenook/.env
ExecStart=/opt/computenook/computenook
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=computenook

[Install]
WantedBy=multi-user.target
