# scripts/hpc-backend.service
[Unit]
Description=HPC Platform Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hpc-platform
EnvironmentFile=-/opt/hpc-platform/.env
ExecStart=/opt/hpc-platform/hpc-backend
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hpc-backend

[Install]
WantedBy=multi-user.target
